# installer/linux/apt/opengeo-suite/debian/prerm
#!/bin/bash
# prerm script for opengeo-suite
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

WEBAPPS=/var/lib/tomcat6/webapps
function remove_apps() {
  invoke-rc.d tomcat6 stop
  for app in $*; do
    if [ -d $WEBAPPS/$app ]; then 
       rm -rf $WEBAPPS/$app
    fi
    if [ -e $WEBAPPS/$app.war ]; then 
       rm -rf $WEBAPPS/$app.war
    fi
  done
  invoke-rc.d tomcat6 start
}

case "$1" in
    remove|deconfigure)

    remove_apps dashboard geoeditor geoexplorer geowebcache recipes 
    ;;

    upgrade|failed-upgrade)
    
    # on upgrade delete all exception geoexplorer, since we want to preserve
    # some its configuration, that is handled postinst
    remove_apps dashboard geoeditor geowebcache recipes 
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
